# .config/systemd/user/sxhkd.service
[Unit]
Description=sxhkd hotkey daemon

[Service]
ExecStart=/usr/bin/sxhkd
Restart=always

[Install]
WantedBy=default.target
